# support/rc.local
#!/bin/bash
sudo sysctl -w net.ipv4.ip_forward=1
/sbin/iptables -P FORWARD ACCEPT
/sbin/iptables --table nat -A POSTROUTING -o eth0 -j MASQUERADE

#For port forwarding :
#sysctl -w net.ipv4.ip_forward=1
#iptables -t nat -A PREROUTING -j DNAT -d ipaddresswanrouter -p tcp --dport (number port) --to (ip address lan server)

#For transparent proxy :
#iptables -t nat -A PREROUTING -i eth1 -p tcp -m tcp --dport 80 -j DNAT --to-destination (ip address lan router):3127

exit 0
